# scripts2/shutdown.service
[Unit]
Description=Shutdown the system

[Service]
User=0
ExecStart=/bin/bash poweroff
