# scripts2/shutdown.timer
[Unit]
Description=Shutdown the system Fridays at 22:00

[Timer]
OnCalendar=Fri *-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
